# services/telegram-scraper/telegram-scraper.service
[Unit]
Description=Telegram insight scraper (local → fund-manager :8787)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/lazywork/workspace/services/telegram-scraper
ExecStart=/home/lazywork/workspace/services/telegram-scraper/.venv/bin/python main.py
Restart=on-failure
RestartSec=30
Environment=TG_INSIGHT_API_URL=http://127.0.0.1:8787/feed/telegram/insight

[Install]
WantedBy=default.target
