# services/telegram-scraper/telegram-scraper.service
[Unit]
Description=Telegram insight scraper (local → fund-manager :8787)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/lazywork/workspace/services
ExecStart=/home/lazywork/workspace/services/telegram-scraper/.venv/bin/python -m telegram_scraper
Restart=on-failure
RestartSec=30
EnvironmentFile=/home/lazywork/workspace/services/telegram-scraper/.env

[Install]
WantedBy=default.target
